# supervisor.conf
[program:boom]
command=/usr/bin/env python services/boom/main.py %(ENV_FLAGS)s
environment=PYTHONPATH=".",PYTHONUNBUFFERED="1"
stdout_logfile=log/boom_service.log
redirect_stderr=true
